# ansible/templates/inventory-manager.service.j2
[Unit]
Description=Inventory Manager Application
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User={{ app_user | default('ec2-user') }}
WorkingDirectory={{ app_dir | default('/opt/inventory-manager') }}
ExecStart=/usr/local/bin/docker-compose up
ExecStop=/usr/local/bin/docker-compose down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
